# meta-noccarc/recipes-extended/network-config/network-config/network-config.service
[Unit]
Description= Network config service V730i

[Service]
Type=idle
ExecStart=/etc/scripts/network-config.sh

[Install]
WantedBy=multi-user.target
