# meta-noccarc/recipes-core/systemd/systemd-v730i-init/V730i_init.service
[Unit]
Description=V510 Initial config service

[Service]
ExecStart=/etc/scripts/V730i_init.sh

[Install]
WantedBy=multi-user.target
